# 009_dashboard/com.violet.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.violet.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/seongwookjang/project/git/violet_sw/009_dashboard/app.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/seongwookjang/project/git/violet_sw/009_dashboard</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/seongwookjang/project/git/violet_sw/009_dashboard/logs/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/seongwookjang/project/git/violet_sw/009_dashboard/logs/dashboard_error.log</string>
</dict>
</plist>
